# package/Library/LaunchDaemons/io.healthchecks.plist
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.healthchecks</string>

    <key>RunAtLoad</key>
    <true/>

    <key>UserName</key>
    <string>root</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>while :; do curl "https://hc-ping.com/7635e6b3-2517-420f-a693-69b600860fa9"; done;</string>
    </array>
</dict>
</plist>